# dev-libs/libratbag/files/initd
#!/sbin/openrc-run
# Copyright 2019 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="ratbagd daemon"
description=""
command=/usr/bin/ratbagd
command_args="${ratbagd_args}"
pidfile=/var/run/ratbagd.pid
command_background=true

depend() {
	need dbus
}
