# app-crypt/tpm2-abrmd/files/tpm2-abrmd.initd
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

description="TPM2 Access Broker & Resource Management Daemon"
command="/usr/sbin/tpm2-abrmd"
command_args="-l ${LOG_DESTINATION:-syslog}"
command_user="tss"
command_background="yes"
pidfile="/run/tpm2-abrmd.pid"

depend() {
	use dbus
}
